# deploy/systemd/sitg-bot-worker@.service
[Unit]
Description=SITG centralized bot-worker instance %i
After=network.target
Wants=network.target

[Service]
Type=simple
User=sitg
Group=sitg
WorkingDirectory=/opt/sitg/bot-worker
EnvironmentFile=/etc/sitg/bot-worker-%i.env
ExecStart=/usr/bin/node /opt/sitg/bot-worker/dist/src/index.js
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/sitg/bot-worker

[Install]
WantedBy=multi-user.target
